# timestamp/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root  . /etc/environment; /usr/local/bin/php /cron_scripts/test_cron.php >> /proc/1/fd/1 2>&1
